# templates/backup.timer
[Unit]
Description=Backup Timer

[Timer]
OnCalendar={{ backup_timer_string }}

[Install]
WantedBy=multi-user.target
